# cron_jobs/app/cron.daily/mark_down_old_users
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.mark_down_old_users