# cron_jobs/app/cron.daily/reminder_expiring_memberships
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.send_reminder_on_expiring_memberships