# cron_jobs/app/cron.daily/remove_old_backups
#!/bin/sh

find /var/backups/postgres* -type f -mtime +14 | xargs rm -Rf
find /var/backups/assets-* -type f -mtime +14 | xargs rm -Rf