# cron_jobs/app/cron.daily/rotate_full_members
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.rotate_user_positions_in_subcategories