# cron_jobs/app/cron.daily/suspend_members_on_expiration
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.suspend_full_members_when_membership_expired