# cron_jobs/db/cron.daily/backup_assets
#!/bin/sh
sudo tar cvfz /usr/local/cft/deploy/capistrano/shared/assets /var/backups/assets-`date +'%Y-%m-%d'`