# cron_jobs/db/cron.daily/copy_backup_to_S3
#!/bin/sh
sudo /home/cftuser/s3sync.rb -r --exclude="\.bak$" /var/backups/ backups.beamazing.co.nz:/