# cron_jobs/app/cron.hourly/count_users
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.count_users