# cron_jobs/app/cron.hourly/process_paid_xero_invoices
#!/bin/sh

ruby /usr/local/cft/deploy/capistrano/current/script/runner -e production TaskUtils.process_paid_xero_invoices